# afo.service.in
[Unit]
Description=Alien Front Online Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=INSTALL_USER
ExecStart=SBINDIR/afoserver SYSCONFDIR/afo.cfg
StandardOutput=append:/var/log/afo.log

[Install]
WantedBy=multi-user.target
